# pion.service
# Pion systemd user service
# Install: cp pion.service ~/.config/systemd/user/
# Enable:  systemctl --user enable pion
# Start:   systemctl --user start pion
# Logs:    journalctl --user -u pion -f

[Unit]
Description=Pion Messaging Bridge
After=network.target

[Service]
Type=simple
# Update these paths to match your setup
WorkingDirectory=%h/Projects/pion
ExecStart=/usr/bin/bun run src/daemon.ts
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
KillMode=control-group
TimeoutStopSec=15
SendSIGKILL=yes
Environment=NODE_ENV=production
Environment=PATH=%h/.local/share/fnm/node-versions/v22.18.0/installation/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
